# deploy/ansible/templates/watchdog/llamenos-watchdog.service.j2
[Unit]
Description=Llamenos Container Watchdog
After=docker.service

[Service]
Type=oneshot
User={{ deploy_user }}
Group={{ deploy_group }}
ExecStart={{ app_dir }}/scripts/container-watchdog.sh
